# dot_config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -c 15 -w 30 -f 0 -m +30 -I ~/.config/dunst/battery-alert.png
